# startup/openrc-init.in
#!/sbin/openrc-run
#
# Copyright (c) 2017 Nagios(R) Core(TM) Development Team
#

command="@sbindir@/ndo2db"
command_args="-c ${NDO2DB_CFG}"
description="Nagios Data Out daemon"
pidfile="@piddir@/ndo2db.pid"

depend() {
    # The Nagios core daemon creates the socket that ndo2db tries to
    # connect to upon starting.
    need mysql nagios
}
